# deploy/runtime/systemd/greatkart-postgres-backup.timer
[Unit]
Description=Run GreatKart PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
